# usr/etc/profile.d/hphp.sh
export USE_HHVM=1
export HPHP_HOME=/usr/lib64/hiphop-php
export HPHP_LIB=/usr/lib64/hiphop-php/bin
export HHVM_LIB_PATH=/usr/lib64/hiphop-php/bin
